# config/systemd/axiom-operator-summary.timer.example
[Unit]
Description=AXIOM Operator Summary Projection Timer (every 15 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
Persistent=true

[Install]
WantedBy=timers.target